# scraper/com.dereks-turo.scrape.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Derek's Turo Metrics — Daily Scraper Launch Agent

  INSTALL:
    cp com.dereks-turo.scrape.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.dereks-turo.scrape.plist

  UNINSTALL:
    launchctl unload ~/Library/LaunchAgents/com.dereks-turo.scrape.plist
    rm ~/Library/LaunchAgents/com.dereks-turo.scrape.plist

  CHECK STATUS:
    launchctl list | grep dereks-turo

  MANUAL RUN:
    cd scraper && npm run scrape

  LOGS:
    tail -f ~/Library/Logs/dereks-turo-scrape.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dereks-turo.scrape</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/joegrady/.nvm/versions/node/v20.20.0/bin/node</string>
    <string>/Users/joegrady/Development/projects/dereks-turo-metrics/scraper/scrape-all.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/joegrady/Development/projects/dereks-turo-metrics/scraper</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/joegrady/Library/Logs/dereks-turo-scrape.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/joegrady/Library/Logs/dereks-turo-scrape.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
